# zendesk-fastapi.service
[Unit]
Description=AWS Zendesk Assistant FastAPI Server
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/aws-zendesk-assistant
ExecStart=/usr/bin/python3 /root/aws-zendesk-assistant/fastapi_server.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
